# lichess_bot.conf
description "Lichess discord bot"
author "thijsvandenhout@live.nl"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/thijs/Lichess-discord-bot/src/lichess_bot.py
